# rsc/etc/profile.d/30-android.sh
export ANDROID_SDK_ROOT="$HOME/.local/share/android-sdk"
export PATH=$ANDROID_SDK_ROOT/cmdline-tools/latest/bin:$ANDROID_SDK_ROOT/platform-tools:$PATH
